# deploy/webster-sec.service
[Unit]
Description=Webster SEC Filing chat
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=webster
WorkingDirectory=/srv/webster-sec-filing
EnvironmentFile=/etc/webster-sec.env
ExecStart=/usr/bin/node server/index.js
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/srv/webster-sec-filing
ProtectHome=true

[Install]
WantedBy=multi-user.target
